# confs/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=5MB
pidfile=/tmp/supervisord.pid
user=wildfly
nodaemon=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:start-wildfly]
command=/start-wildfly.sh
priority=1
user=wildfly
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_events_enabled=true

# [program:jboss-cli]
# command=/opt/jboss/wildfly/bin/jboss-cli.sh --file=/opt/jboss/wildfly/bin/create_infinity_datasources.cli 
# priority=999
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stdout_events_enabled=true
# stderr_events_enabled=true
